# cron.weekly
#! /bin/bash

set -euo pipefail

suffix=geograph-update-bot
srcdir=/home/bjh21/geograph-update-bot
workdir="/quicksand/bjharris.${suffix}"

quicksand-setup --rmdir "${workdir}" || true
quicksand-setup --mkdir "${workdir}"

cp -a "${srcdir}" "${workdir}"
cd "${workdir}"

./mkvenv venv

( cd geograph-db && ./geograph_fetch && ./geograph_import )

# XXX actually run bot here.
